# asis.conf
description "uWSGI server instance configured to serve ASIS"

start on runlevel [2345]
stop on runlevel [!2345]

setuid vagrant
setgid www-data

chdir /workspace/ASIS1
exec uwsgi --ini asis.ini